# fs/etc/s6-overlay/s6-rc.d/install-monitoring-key/up
#!/command/with-contenv sh
install-monitoring-key
