# fs/etc/s6-overlay/s6-rc.d/initialize-openssh/up
#!/command/with-contenv sh
initialize-openssh
