# fs/etc/s6-overlay/s6-rc.d/configure-proxy/up
#!/command/with-contenv sh
configure-proxy
